# dcec-autopilot.service
[Unit]
Description=dcec-autopilot

[Service]
ExecStart=/dcec-autopilot.sh

[Install]
WantedBy=multi-user.target